# contrib/systemd/traefik.service
[Unit]
Description=Traefik

[Service]
ExecStart=/usr/bin/traefik --configFile=/etc/traefik.toml
Restart=on-failure
