# contrib/packaging/deb/postinst
#!/bin/bash

# skip systemd if not available (for example in WSL)
if pidof systemd > /dev/null; then
    systemctl daemon-reload
    systemctl enable traefik
    systemctl start traefik
fi

